# botcha.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 2

setuid ubuntu
cd ~/botcha

env PORT=3010

exec node Index.js